# packer/rproxy/resources/crons.cron
*/10 * * * * /usr/local/bin/update-kube-nodes.sh
0 4 * * * certbot renew
