# packer/rproxy/resources/update-kube-nodes.cron
* * * * * /usr/local/bin/update-kube-nodes.sh
